# docker/ros_melodic_aws/scripts/crontab
*/5 * * * * user sudo logrotate /etc/logrotate.d/ros2_logrotate
# Empty line